# .config/autostart/google-chrome.desktop
[Desktop Entry]
Type=Application
Terminal=false
Exec=/opt/google/chrome/google-chrome --no-startup-window
Name=Google Chrome
